# system/udev/rules.d/99-fix-ghostkey.rules
ACTION=="add", SUBSYSTEM=="serio", DRIVERS=="atkbd", RUN+="/usr/bin/setkeycodes 56 255"
